# etc/profile.d/git.sh
#!/usr/bin/env bash

export GIT_HOME=/usr/local/lib/git
export PATH=${GIT_HOME}/bin:${PATH}
