# etc/profile.d/nodejs.sh
#!/usr/bin/env bash

export NODEJS_HOME=/usr/local/lib/nodejs
export PHANTOMJS_HOME=/usr/local/lib/phantomjs
export CASPERJS_HOME=/usr/local/lib/casperjs

export PATH=${NODEJS_HOME}:${PHANTOMJS_HOME}/bin:${CASPERJS_HOME}/bin:${PATH}

