# etc/profile.d/scala.sh
#!/usr/bin/env bash

export JAVA_HOME=/usr/local/lib/java
export SBT_HOME=/usr/local/lib/sbt

export PATH=${PATH}:${JAVA_HOME}/bin:${SBT_HOME}/bin
